# deploy/happy-star.service
[Unit]
Description=Happy Star family points
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/wh/apps/happy-star
Environment=PORT=8080
Environment=HOST=0.0.0.0
Environment=HAPPY_STAR_DATA=/home/wh/apps/happy-star/data
# If node is installed via nvm, replace /usr/bin/env node with the absolute
# path from `which node` (systemd does not load the user's nvm PATH).
ExecStart=/usr/bin/env node server/src/index.js
Restart=always
User=wh

[Install]
WantedBy=multi-user.target
